# labs/ssh-tunnel/hostc/_system/etc/rc.local
#!/bin/bash

route add -net 172.50.0.0 netmask 255.255.255.0 gw 172.60.0.76

route add -net 172.70.0.0 netmask 255.255.255.0 gw 172.60.0.75
