# scripts/designer/system/etc/systemd/system/mynotify.service
[Unit]
Description=mynotify

[Service]
ExecStartPre=/usr/bin/wait_param.sh
ExecStart=/sbin/mynotify.py

[Install]
WantedBy=multi-user.target
