# labs/dmz-example/dns/_system/etc/rc.local
#!/bin/bash
/usr/bin/set_default_gw.sh 198.18.1.2
route del -host 172.17.0.1

echo "198.18.1.2" > /etc/resolv.conf
