# labs/dns/gw/_system/etc/rc.local
#!/bin/bash
#

#
# restart the DNS forwarder (rely on DNS provided by ISP)
#
echo "nameserver 192.168.10.2" >/etc/resolv.conf
/etc/init.d/dnsmasq restart

