# labs/routing-basics/gateway/_system/etc/rc.local
/usr/bin/set_default_gw.sh 192.170.0.101

