# labs/backups2/router2/_system/etc/rc.local
#!/bin/bash
route delete default

route add -host 172.25.0.2 gw 172.25.2.3
