# labs/vpnlab2/server/_system/lib/systemd/system/httpserver.service
[Unit]
Description=HTTP Service

[Service]
Type=simple
ExecStart=/home/ubuntu/MyHTTPServer.py
StandardOutput=null

[Install]
WantedBy=multi-user.target
Alias=httpserver.service
