# labs/iptables2/server/_system/etc/xinetd.d/telnetx
service telnetx
{
    flags = REUSE
    socket_type = stream
    port = 2023
    wait = no
    user = root
    server = /usr/sbin/in.telnetd
    log_on_failure += USERID
    disable = no
}
