# labs/dmz-lab/ws2/_system/etc/rc.local
#!/bin/bash
# Use gateway as the name server
#
echo "nameserver 198.18.1.3" > /etc/resolv.conf
/usr/bin/set_default_gw.sh 198.18.1.INNER_LAN2
# remove the dns route so as to not confuse student
sudo route del -host 172.17.0.1
