# scripts/designer/workspace/system/etc/systemd/system/labdump.service
[Unit]
Description=labdump

[Service]
ExecStart=/bin/labdump_server.py

[Install]
WantedBy=multi-user.target
