# scripts/designer/workspace/system/etc/systemd/system/labdump_client.service
[Unit]
Description=labdump_client

[Service]
ExecStart=/bin/start_labdump.sh

[Install]
WantedBy=multi-user.target
