# labs/dhcp-test/client/_system/etc/rc.local
#!/bin/bash
/sbin/dhclient-labtainer eth0
