# labs/dhcp-test/client/_system/etc/rc.local
#!/bin/bash
echo "dl $DISTRIBUTED_LABTAINER" >> /tmp/rc.log
env >> /tmp/rc.log
if [ "$DISTRIBUTED_LABTAINER" == client ]; then
    /sbin/dhclient-labtainer eth0
fi
