# labs/backups2/router1/_system/etc/rc.local
#!/bin/bash
route delete default

route add -host 172.25.3.3 gw 172.25.2.4
