# labs/plc-forensics-adv/plc/_system/lib/systemd/system/httpserver.service
[Unit]
Description=HTTP Service

[Service]
Type=simple
WorkingDirectory=/var/www
ExecStart=/sbin/MyHTTPServer.py
StandardOutput=null

[Install]
WantedBy=multi-user.target
Alias=httpserver.service
