# labs/remote-dns/attacker/_system/etc/rc.local
#!/bin/bash
echo "nameserver 203.0.113.10" > /etc/resolv.conf
