# host_scripts/testscripts/smoke.cron
0 2 * * * /home/mike/SEED/host_scripts/testscripts/powersmoke.sh
