# labs/iptables2/server/_system/usr/lib/systemd/system/wizbang.service
[Unit]
Description=wizbang

[Service]
ExecStart=/home/ubuntu/server

[Install]
WantedBy=multi-user.target
Alias=wizbang.service
