# labs/snort/snort/_system/etc/rc.local
/usr/bin/set_default_gw.sh 192.168.3.10
sudo sysctl -w net.ipv4.ip_forward=0

