# labs/snort/snort/_system/etc/rc.local
#!/bin/bash
/usr/bin/set_default_gw.sh 192.168.3.10
sudo sysctl -w net.ipv4.ip_forward=0

echo "nameserver 192.168.3.10" > /etc/resolv.conf
route del -host 172.17.0.1

