# labs/arp-spoof/webserver/_system/etc/rc.local
#!/bin/bash
/usr/bin/set_default_gw.sh 172.35.0.2
