# labs/ssh-tunnel/hostb/_system/etc/rc.local
#!/bin/bash

route add -net 172.40.0.0 netmask 255.255.255.0 gw 172.50.0.76

route add -net 172.60.0.0 netmask 255.255.255.0 gw 172.50.0.75
