# labs/snort/ws2/_system/etc/rc.local
/usr/bin/set_default_gw.sh 192.168.2.10
route delete default
route add default gw 192.168.2.10
