# labs/snort/ws2/_system/etc/rc.local
#!/bin/bash
/usr/bin/set_default_gw.sh 192.168.2.10
route delete default
route add default gw 192.168.2.10
echo "nameserver 192.168.2.10" > /etc/resolv.conf
route del -host 172.17.0.1

