# labs/local-dns/apollo-dns/_system/etc/rc.local
#!/bin/bash
/usr/bin/set_default_gw.sh 192.168.0.1
echo "192.168.0.1" > /etc/resolv.conf
