# labs/usbtest/other/_system/etc/udev/rules.d/10-usbmount.rules
KERNEL=="sd*", RUN+="/home/ubuntu/iadrive/doiadrive.sh -u"
