# labs/nmaplab/server/_system/etc/xinetd.d/ssh
service ssh
{
          socket_type            = stream
          protocol               = tcp
          wait                   = no
          port                   = 22
          user                   = root
          server                 = /usr/sbin/sshd
          server_args            = -i
          disable                = yes
}

