# labs/softplc2/slave-device/_system/etc/rc.local
#!/bin/bash
export DISPLAY=:0
/home/ubuntu/server.py &
