# labs/ssl/plc2/_system/etc/systemd/system/mynotify.service
[Unit]
Description=mynotify

[Service]
ExecStart=/sbin/mynotify.py

[Install]
WantedBy=multi-user.target
