# labs/dns/isp/_system/etc/rc.local
#!/bin/bash
/etc/init.d/dnsmasq restart
