# labs/dmz-lab/ws1/_system/etc/rc.local
#!/bin/bash
# Use gateway as the name server
#
date >/tmp/rc.log
echo "nameserver 198.18.1.3" > /etc/resolv.conf
/usr/bin/set_default_gw.sh 198.18.1.INNER_LAN1
# remove the dns route so as to not confuse student
route del -host 172.17.0.1
