# labs/snort/remote_ws/_system/etc/rc.local
#/usr/bin/set_default_gw.sh 172.16.0.10
route delete default
route add default gw 172.16.0.10

