# labs/cgc/server/_system/etc/xinetd.d/audio_visualizer
service audio_visualizer
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Audio_Visualizer
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}