# labs/cgc/server/_system/etc/xinetd.d/barcoder
service barcoder
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Barcoder
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}