# labs/cgc/server/_system/etc/xinetd.d/cablegrindllama
service cablegrindllama
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/CableGrindLlama
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}