# labs/cgc/server/_system/etc/xinetd.d/casino_games
service casino_games
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Casino_Games
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}