# labs/cgc/server/_system/etc/xinetd.d/childs_game
service childs_game
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Childs_Game
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}