# labs/cgc/server/_system/etc/xinetd.d/diophantine_password_wallet
service diophantine_password_wallet
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Diophantine_Password_Wallet
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}