# labs/cgc/server/_system/etc/xinetd.d/enslavednode_chat
service enslavednode_chat
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Enslavednode_chat
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}