# labs/cgc/server/_system/etc/xinetd.d/finicky_file_folder
service finicky_file_folder
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Finicky_File_Folder
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}