# labs/cgc/server/_system/etc/xinetd.d/grit
service grit
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Grit
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}