# labs/cgc/server/_system/etc/xinetd.d/kaprica_script_interpreter
service kaprica_script_interpreter
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Kaprica_Script_Interpreter
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}