# labs/cgc/server/_system/etc/xinetd.d/kty_pretty_printer
service kty_pretty_printer
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/KTY_Pretty_Printer
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}