# labs/cgc/server/_system/etc/xinetd.d/material_temperature_simulation
service material_temperature_simulation
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Material_Temperature_Simulation
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}