# labs/cgc/server/_system/etc/xinetd.d/mathematical_solver
service mathematical_solver
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Mathematical_Solver
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}