# labs/cgc/server/_system/etc/xinetd.d/matrix_math_calculator
service matrix_math_calculator
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Matrix_Math_Calculator
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}