# labs/cgc/server/_system/etc/xinetd.d/matrices_for_sale
service matrices_for_sale
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/matrices_for_sale
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}