# labs/cgc/server/_system/etc/xinetd.d/message_service
service message_service
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Message_Service
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}