# labs/cgc/server/_system/etc/xinetd.d/movie_rental_service
service movie_rental_service
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Movie_Rental_Service
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}