# labs/cgc/server/_system/etc/xinetd.d/multicast_chat_server
service multicast_chat_server
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Multicast_Chat_Server
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}