# labs/cgc/server/_system/etc/xinetd.d/multi_arena_pursuit_simulator
service multi_arena_pursuit_simulator
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Multi_Arena_Pursuit_Simulator
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}