# labs/cgc/server/_system/etc/xinetd.d/multipass2
service multipass2
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Multipass2
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}