# labs/cgc/server/_system/etc/xinetd.d/on_sale
service on_sale
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/On_Sale
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}