# labs/cgc/server/_system/etc/xinetd.d/one_amp
service one_amp
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/One_Amp
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}