# labs/cgc/server/_system/etc/xinetd.d/online_job_application2
service online_job_application2
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/online_job_application2
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}