# labs/cgc/server/_system/etc/xinetd.d/packet_analyzer
service packet_analyzer
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Packet_Analyzer
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}