# labs/cgc/server/_system/etc/xinetd.d/packet_receiver
service packet_receiver
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Packet_Receiver
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}