# labs/cgc/server/_system/etc/xinetd.d/pattern_finder
service pattern_finder
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Pattern_Finder
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}