# labs/cgc/server/_system/etc/xinetd.d/ram_based_filesystem
service ram_based_filesystem
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/RAM_based_filesystem
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}