# labs/cgc/server/_system/etc/xinetd.d/resort_modeller
service resort_modeller
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Resort_Modeller
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}