# labs/cgc/server/_system/etc/xinetd.d/sad_face_template_engine_sfte
service sad_face_template_engine_sfte
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Sad_Face_Template_Engine_SFTE
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}