# labs/cgc/server/_system/etc/xinetd.d/sauth
service sauth
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/SAuth
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}