# labs/cgc/server/_system/etc/xinetd.d/secure_compression
service secure_compression
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Secure_Compression
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}