# labs/cgc/server/_system/etc/xinetd.d/simple_integer_calculator
service simple_integer_calculator
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/simple_integer_calculator
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}