# labs/cgc/server/_system/etc/xinetd.d/spiffs
service spiffs
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/SPIFFS
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}