# labs/cgc/server/_system/etc/xinetd.d/stack_vm
service stack_vm
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/stack_vm
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}