# labs/cgc/server/_system/etc/xinetd.d/stream_vm2
service stream_vm2
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/stream_vm2
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}