# labs/cgc/server/_system/etc/xinetd.d/string_storage_and_retrieval
service string_storage_and_retrieval
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/String_Storage_and_Retrieval
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}