# labs/cgc/server/_system/etc/xinetd.d/string_info_calculator
service string_info_calculator
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/String_Info_Calculator
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}