# labs/cgc/server/_system/etc/xinetd.d/the_longest_road
service the_longest_road
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/The_Longest_Road
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}