# labs/cgc/server/_system/etc/xinetd.d/thermal_controller_v2
service thermal_controller_v2
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Thermal_Controller_v2
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}