# labs/cgc/server/_system/etc/xinetd.d/tiaca
service tiaca
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/TIACA
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}