# labs/cgc/server/_system/etc/xinetd.d/tick-a-tack
service tick-a-tack
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Tick-A-Tack
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}