# labs/cgc/server/_system/etc/xinetd.d/tvs
service tvs
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/TVS
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}