# labs/cgc/server/_system/etc/xinetd.d/university_enrollment
service university_enrollment
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/university_enrollment
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}