# labs/cgc/server/_system/etc/xinetd.d/user_manager
service user_manager
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/User_Manager
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}