# labs/cgc/server/_system/etc/xinetd.d/vector_graphics_format
service vector_graphics_format
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Vector_Graphics_Format
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}