# labs/cgc/server/_system/etc/xinetd.d/vector_graphics_2
service vector_graphics_2
{
    flags = REUSE
    socket_type = stream
    wait = no
    user = root
    server = /usr/sbin/Vector_Graphics_2
    log_on_failure += USERID
    log_on_success += USERID
    disable = no
}