# labs/routing-basics/remote_ws/_system/etc/rc.local
#!/bin/bash
/usr/bin/set_default_gw.sh 172.16.0.10
echo "nameserver 172.16.0.10" > /etc/resolv.conf
route del -host 172.17.0.1

