# labs/routing-basics/isp/_system/etc/rc.local
#!/bin/bash
lan1=$(ifconfig | grep -B1 "inet.*203.0.113.1" | awk '$1!="inet" && $1!="--" {print substr($1,1,length($1)-1)}')
wan=$(ifconfig | grep -B1 "inet.*10.10.0.1" | awk '$1!="inet" && $1!="--" {print substr($1,1,length($1)-1)}')
sudo iptables --table nat -I POSTROUTING 1 --out-interface $wan -j MASQUERADE
#sudo iptables --append FORWARD --in-interface $lan1 -j ACCEPT

systemctl restart dnsmasq
